# debian/postinst
#!/bin/sh
update-rc.d -f hostapd defaults
